# systemd/system/bless-boot.service
[Unit]
Description=Bless Current Boot - Shani OS Boot Counting
After=multi-user.target data.mount mark-boot-success.service
Wants=data.mount
# Only bless the boot after mark-boot-success has confirmed the system reached
# multi-user.target and written /data/boot-ok. This ensures the systemd-boot
# boot counter (+3-0 tries) is only satisfied when the system is truly healthy,
# not at sysinit time when nothing useful has started yet.
ConditionPathExists=/data/boot-ok

[Service]
Type=oneshot
ExecStart=bootctl set-good
SuccessExitStatus=0 1
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
