# systemd/system/shani-user-setup.service
[Unit]
Description=Sync groups, shell, and skel for all regular users
# Triggered by shani-user-setup.path when:
#   - /data/overlay/etc/upper/passwd is created or changed (user added at runtime)
#   - /data/user-setup-needed exists (first install or slot switch via shani-deploy)
#   - /etc/skel is modified (skel update via slot switch)
#
# /etc is mounted by dracut pre-pivot hook (shanios-overlay-etc.sh) before
# systemd PID 1 starts. No etc.mount unit exists — the overlay is live at boot.
# ConditionPathIsMountPoint=/etc guards against the rare case the overlay failed.
#
# Condition* with | prefix = OR logic among | conditions (AND with non-| ones).
# Service runs if:
#   (/data/user-setup-needed exists OR /data/overlay/etc/upper/passwd exists)
#   AND /etc is a mountpoint
# This covers both triggers:
#   - deploy marker: slot switch / first install (shani-deploy writes the marker)
#   - passwd upper layer: user added at runtime (useradd copies up /etc/passwd)
ConditionPathExists=|/data/user-setup-needed
ConditionPathExists=|/data/overlay/etc/upper/passwd
ConditionPathIsMountPoint=/etc

# Wants= (not Requires=) so a missing accounts-daemon does not block user setup.
Wants=accounts-daemon.service
After=accounts-daemon.service data.mount
# data.mount is a hard requirement — needed to read/clear the user-setup-needed marker.
Requires=data.mount

[Service]
Type=oneshot
ExecStart=/usr/local/bin/shani-user-setup
# Clear the deploy marker after a successful run so the service does not
# re-trigger on the next boot. The .path unit will still re-trigger it if a
# new user is added (via passwd change) or shani-deploy writes the marker again.
ExecStartPost=-/bin/rm -f /data/user-setup-needed
TimeoutStartSec=120s

# No [Install] — activated exclusively by shani-user-setup.path.
