# systemd/user/flatpak-update-user.timer
# Install to: /etc/systemd/user/flatpak-update-user.timer
# (NOT /etc/systemd/system/ — this is a user-scope timer)
[Unit]
Description=Periodic Flatpak Update (User)

[Timer]
OnBootSec=20min
OnUnitActiveSec=12h
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
