# systemd/user/shani-update.service
[Unit]
Description=Shani OS Update Check UI
After=graphical-session.target network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# The script auto-detects DISPLAY; override here if necessary.
# Environment=DISPLAY=:0
ExecStart=/usr/local/bin/shani-update
TimeoutSec=300

[Install]
WantedBy=graphical.target

